# scripts/groqtype.service.example
# Example systemd unit for GroqType (generated automatically by install.sh).
# Do not edit this file directly — run ./scripts/install.sh or ./scripts/doctor.sh --fix
# to install/regenerate the live unit with correct paths and session environment.

[Unit]
Description=GroqType speech-to-text daemon
After=network.target sound.target ydotool.service keyd.service graphical.target
Wants=ydotool.service keyd.service

[Service]
Type=simple
ExecStart=/path/to/GroqType/.venv/bin/python /path/to/GroqType/groqtype.py daemon
Restart=on-failure
RestartSec=2
WorkingDirectory=/path/to/GroqType
User=youruser
Group=youruser
SupplementaryGroups=keyd input
Environment=GROQTYPE_CONFIG=/etc/groqtype/config.json
Environment=PYTHONPATH=/path/to/GroqType
Environment=GROQ_API_KEY=
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=YDOTOOL_SOCKET=/run/user/1000/.ydotool_socket
Environment=PATH=/home/youruser/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
